# server-stuff/mount-server-rc
#!/sbin/openrc-run
###depend() {
###    need localmount sysfs
###    before lxc.files
###}

start(){
bash /etc/init.d/mount-server start
}

stop(){
bash /etc/init.d/mount-server stop
}

